# supervisord_rq.conf
[supervisord]
nodaemon=true

[program:rq_worker]
command=/bin/bash /backend/start_rq_worker.sh

[program:mlflow_ui]
command=/bin/bash /backend/start_mlflow_ui.sh